# environments/uwsgi/flower.conf
[program:flower]
directory = /app/
command = /usr/local/bin/celery flower -A celery_worker.celery --loglevel=info